# templates/backup.cron
0 */12 * * * /backup.sh >> /var/log/backup.log 2>&1